# run/systemd/system/run-11179.service.d/50-Description.conf
[Unit]
Description=yum locked retry of update-motd
