# run/systemd/system/run-11179.service.d/50-ExecStart.conf
[Service]
ExecStart=
ExecStart=@/usr/sbin/update-motd "/usr/sbin/update-motd"
