# run/systemd/generator.late/cfn-hup.service
# Automatically generated by systemd-sysv-generator

[Unit]
Documentation=man:systemd-sysv-generator(8)
SourcePath=/etc/rc.d/init.d/cfn-hup
Description=SYSV: Runs user-specified actions when a

[Service]
Type=forking
Restart=no
TimeoutSec=5min
IgnoreSIGPIPE=no
KillMode=process
GuessMainPID=no
RemainAfterExit=no
PIDFile=/var/run/cfn-hup.pid
ExecStart=/etc/rc.d/init.d/cfn-hup start
ExecStop=/etc/rc.d/init.d/cfn-hup stop
ExecReload=/etc/rc.d/init.d/cfn-hup reload
